# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live
nodaemon=true                                   ; run on foreground

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock           ; use a unix:// URL  for a unix socket

[program:nginx]
command=nginx
stderr_logfile = /var/log/supervisord/nginx-stderr.log
stdout_logfile = /var/log/supervisord/nginx-stdout.log

[program:goofys]
command=goofys -f vod.eiplus.tv /stream/s3_vod/
stderr_logfile = /var/log/supervisord/goofys-stderr.log
stdout_logfile = /var/log/supervisord/goofys-stdout.log

[program:api]
command=npm start
directory=/stream/api
stderr_logfile = /var/log/supervisord/api-stderr.log
stdout_logfile = /var/log/supervisord/api-stdout.log